# app-notify.service
[Unit]
Description=Notify Service
After=network.target redis.service mongod.service
Wants=redis.service mongod.service

[Service]
Type=simple
User=dev
Group=dev
WorkingDirectory=/opt/notify
ExecStart=/opt/notify/notify
EnvironmentFile=/opt/notify/.env
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5

# Graceful shutdown: send SIGTERM, wait 15s, then SIGKILL
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=15

# Resource limits
LimitNOFILE=65535
MemoryMax=512M
CPUQuota=200%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/notify

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=app-notify

[Install]
WantedBy=multi-user.target
